# release/packaging/systemd/doublezero-device-health-oracle.service
# NOTE: This is a placeholder unit.
# It will NOT start the device health oracle with correct runtime flags.
# Expected: Ansible will manage an override file at:
#   /etc/systemd/system/doublezero-device-health-oracle.service.d/override.conf
# which will set the real ExecStart and options.

[Unit]
Description=DoubleZero Device Health Oracle
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/doublezero-device-health-oracle
Restart=on-failure

[Install]
WantedBy=multi-user.target
