# release/packaging/systemd/doublezero-activator.service
# NOTE: This is a placeholder unit.
# It will NOT start the activator with correct runtime flags.
# Expected: Ansible will manage an override file at:
#   /etc/systemd/system/doublezero-activator.service.d/override.conf
# which will set the real ExecStart and options.

[Unit]
Description=DoubleZero Activator
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/doublezero-activator --onchain-allocation
Restart=on-failure

[Install]
WantedBy=multi-user.target
