# release/packaging/systemd/doublezero-funder.service
# NOTE: This is a placeholder unit.
# It will NOT start the funder with correct runtime flags.
# Expected: Ansible will manage an override file at:
#   /etc/systemd/system/doublezero-funder.service.d/override.conf
# which will set the real ExecStart and options.

[Unit]
Description=DoubleZero Funder
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/doublezero-funder
Restart=on-failure

[Install]
WantedBy=multi-user.target
