# release/packaging/systemd/doublezero-controller.service
# NOTE: This is a placeholder unit.
# It will NOT start the controller with correct runtime flags.
# Expected: Ansible will manage an override file at:
#   /etc/systemd/system/doublezero-controller.service.d/override.conf
# which will set the real ExecStart and options.

[Unit]
Description=DoubleZero Controller
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/doublezero-controller
Restart=on-failure

[Install]
WantedBy=multi-user.target
