# controlplane/internet-latency-collector/cmd/collector/doublezero-internet-latency-collector.service
[Unit]
Description=DoubleZero internet latency collector
After=network.target

[Service]
Type=simple
EnvironmentFile=/etc/doublezero-internet-latency-collector/config.env
ExecStart=/usr/local/bin/doublezero-internet-latency-collector run --env testnet --keypair /etc/doublezero-internet-latency-collector/keypair.json
Restart=always
RestartSec=10
User=dz-internet-latency
Group=dz-internet-latency

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/doublezero-internet-latency-collector

[Install]
WantedBy=multi-user.target
